# deploy/minierp.service
[Unit]
Description=MiniErp Dockerized (web + celery)
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=abc
Group=docker
WorkingDirectory=/home/abc/MiniErp
ExecStart=/usr/bin/bash -c 'docker start minierp-web minierp-celery 2>/dev/null || true'
ExecStop=/usr/bin/bash -c 'docker stop minierp-web minierp-celery 2>/dev/null || true'
ExecReload=/usr/bin/bash -c 'docker restart minierp-web minierp-celery 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
