# deploy/minierp.service
[Unit]
Description=MiniErp Project Management System
After=network.target nginx.service

[Service]
Type=simple
User=abc
Group=abc
WorkingDirectory=/home/abc/MiniErp
Environment="PATH=/home/abc/MiniErp/venv/bin"
Environment="PYTHONPATH=/home/abc/MiniErp"
Environment="DJANGO_SETTINGS_MODULE=minierp.settings"
ExecStart=/home/abc/MiniErp/venv/bin/gunicorn --config /home/abc/MiniErp/deploy/gunicorn_config.py minierp.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target